# server/supervisord.conf
; Prod supervisor config file.

[supervisord]
logfile=/tmp/supervisord.log ; main log file; default $CWD/supervisord.log
pidfile=/tmp/supervisord.pid ; supervisord pidfile; default supervisord.pid

[program:gunicorn]
environment=PRIMES_ENV=production
; TODO: use deploy user
command=/root/notre-dame-prime/server/env/bin/gunicorn 'server:load_app("production")' -c ./gunicorn.conf.py ;
directory=/root/notre-dame-prime/server ;
user=deploy ;
autostart=true ;
autorestart=true ;
redirect_stderr=true ;
