# config/supervisord.conf
; Prod supervisor config file.

[unix_http_server]
file=/tmp/supervisor.sock ;

[supervisord]
logfile=/tmp/supervisord.log ; main log file; default $CWD/supervisord.log
pidfile=/tmp/supervisord.pid ; supervisord pidfile; default supervisord.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ;

[program:gunicorn]
; TODO: use deploy user
environment=GOOGLE_APPLICATION_CREDENTIALS='/root/paint-by-primes/server/resources/serviceAccount.json'
; TODO: use deploy user
command=/root/paint-by-primes/server/env/bin/gunicorn 'server:load_app("production")' -c ../config/gunicorn.conf.py ;
directory=/root/paint-by-primes/server ;
; TODO: use deploy user
user=root ;
autostart=true ;
autorestart=true ;
redirect_stderr=true ;
